# hassio-addon-nest-protect/rootfs/etc/services.d/nest_proetct/run
#!/command/with-contenv bashio
# ==============================================================================

bashio::log.info "Starting nest protect..."

exec /usr/src/app/run.sh
