# hassio-addon-enedis-power-cut/rootfs/etc/services.d/enedis_power_cut/run
#!/command/with-contenv bashio
# ==============================================================================

bashio::log.info "Starting Enedis Power Cut..."

exec /usr/src/app/run.sh
